# gunicorn/yats.conf
[program:yats]
directory=/var/web/yats/web
command=GUNICORN_BIN web.wsgi -c /var/web/yats/gunicorn/gunicorn.conf